# app-misc/kmonad/kmonad.init
#!/sbin/openrc-run

pidfile="/run/kmonad.pid"
command="/usr/bin/kmonad"
command_args="/etc/conf.d/kmonad/config.kbd"
command_background=1

depend() {
	after coldplug
	need localmount
}
